# systemd/alice-lg.service
[Unit]
Description=Alice Looking Glass - BGP Looking Glass Frontend
Documentation=https://github.com/alice-lg/alice-lg
After=network-online.target birdwatcher.service
Requires=birdwatcher.service

[Service]
Type=simple
User=alice-lg
Group=alice-lg
WorkingDirectory=/var/lib/alice-lg
ExecStart=/usr/local/bin/alice-lg -config /etc/alice-lg/alice.conf
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/alice-lg

[Install]
WantedBy=multi-user.target
